# time-fixed/services/minecraft-announce.timer
[Unit]
Description=Timer for minecraft-announce

[Timer]
# 毎日01:55に実行
OnCalendar=*-*-* 01:55:00

[Install]
WantedBy=timers.target
