# time-fixed/services/minecraft-stop.timer
[Unit]
Description=Timer for minecraft-stop

[Timer]
# 毎日02:00に実行
OnCalendar=*-*-* 02:00:00

[Install]
WantedBy=timers.target
